# systemd/vortexl2-forward-daemon.service
[Unit]
Description=VTP Forwarding Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=PYTHONPATH=/opt/vortexl2
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/usr/bin/python3 -m vortexl2.forward_daemon
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
